# suricata/cron.d/suricata
0 18 * * * root /usr/bin/suricata-update  >/dev/null 2>&1 >/dev/stdout